# backup/systemd/backup.timer
[Unit]
Descriptions=rclone2cloud

[Timer]
Unit=backup.service
OnBootSec=1min
OnUnitInactiveSec=5min

[Install]
WantedBy=timers.target
